# package/dweb/root/etc/uci-defaults/20_luci_operation
#!/bin/sh
uci set network.op=operation
uci set network.op.mode=gateway
uci commit
