# package/dweb/root/etc/uci-defaults/30_3g_proto
#!/bin/sh
uci set network.wan2=interface
uci set network.wan2.ifname="ppp0"
uci set network.wan2.proto="3g"
uci set network.wan2.enabled=0
uci set network.wan2.metric=20
uci commit
