# package/dweb/root/etc/uci-defaults/50_theme_dweb
#!/bin/sh
uci batch <<-EOF
	set luci.themes.dweb=/luci-static/dweb
	set luci.main.mediaurlbase=/luci-static/dweb
	commit luci
EOF
exit 0
